# packaging/deb/debian/postinst
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload
  systemctl enable neon-fix-drm.path 2>/dev/null || true
fi
